# rev_tun/templates/supervisor.conf.j2
[program:{{ name }}]
command={{ command }}
autostart=true
autorestart=true
startretries={{ retry }}
stderr_logfile={{ log_dir }}/{{ name }}.err.log
stdout_logfile={{ log_dir }}/{{ name }}.out.log
